# supervisord.conf
; =============================================================================
; Supervisord Configuration — manages ComfyUI + FastAPI in a single container
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

; =============================================================================
; ComfyUI Instance 0 — GPU 0 on port 8188
; =============================================================================
[program:comfyui_0]
command=bash /workspace/app-root/scripts/run_comfyui.sh --listen 127.0.0.1 --port 8188 --dont-print-server --user-directory /workspace/comfyui_user_0 --database-url sqlite:///:memory: %(ENV_COMFYUI_ARGS)s
directory=/workspace/ComfyUI
autostart=true
autorestart=true
startretries=5
startsecs=15
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=CUDA_VISIBLE_DEVICES="0"

; =============================================================================
; ComfyUI Instance 1 — GPU 1 on port 8189
; =============================================================================
[program:comfyui_1]
command=bash /workspace/app-root/scripts/run_comfyui.sh --listen 127.0.0.1 --port 8189 --dont-print-server --user-directory /workspace/comfyui_user_1 --database-url sqlite:///:memory: %(ENV_COMFYUI_ARGS)s
directory=/workspace/ComfyUI
autostart=false
autorestart=true
startretries=5
startsecs=15
priority=11
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=CUDA_VISIBLE_DEVICES="1"

; =============================================================================
; FastAPI — Main API gateway on port 45701 (accesses all GPUs)
; =============================================================================
[program:fastapi]
command=python -m uvicorn app.main:app --host 0.0.0.0 --port 45701 --workers 1 --log-level info
directory=/workspace/app-root
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=CUDA_VISIBLE_DEVICES="0,1",PYTHONPATH="/workspace/app-root"
