# supervisord.conf
; =============================================================================
; Supervisord Configuration — manages ComfyUI + FastAPI in a single container
; =============================================================================

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

; =============================================================================
; ComfyUI — Headless SDXL server on port 8188
; =============================================================================
[program:comfyui]
command=python /workspace/ComfyUI/main.py --listen 127.0.0.1 --port 8188 --dont-print-server
directory=/workspace/ComfyUI
autostart=true
autorestart=true
startretries=5
startsecs=15
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=CUDA_VISIBLE_DEVICES="0"

; =============================================================================
; FastAPI — Main API gateway on port 8000
; =============================================================================
[program:fastapi]
command=python -m uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1 --log-level info
directory=/workspace/app-root
autostart=true
autorestart=true
startretries=5
startsecs=10
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=CUDA_VISIBLE_DEVICES="0",PYTHONPATH="/workspace/app-root"
